# .etc/udev/rules.d/mouse.rules
ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c53f", RUN+="~/mouse.sh"
